# dist/nethvoice-report-api.service
[Unit]
Description=NethVoice Report API unit
After=syslog.target network-online.target systemd-tmpfiles-setup.service
Documentation=https://github.com/nethesis/nethvoice-report

[Service]
EnvironmentFile=/opt/nethvoice-report/api/api.env
ExecStart=/opt/nethvoice-report/api/api
User=asterisk

[Install]
WantedBy=multi-user.target