# root/etc/cron.d/nethvoice-report-tasks
# precalculate queries during night
59 23 * * *	asterisk /opt/nethvoice-report/tasks/tasks queries
