# root/etc/cron.d/nethvoice-report-tasks
# precalculate phonebook numbers during night
00 00 * * *     asterisk /opt/nethvoice-report/tasks/tasks phonebook

# precalculate filter values during night
10 00 * * *     asterisk /opt/nethvoice-report/tasks/tasks values

# precalculate queries during night
20 00 * * *	asterisk /opt/nethvoice-report/tasks/tasks queries
